# scripts/launchd/com.agentic-sniper.paper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic-sniper.paper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>-lc</string>
        <string>export PATH="$HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:$PATH" &amp;&amp; /usr/bin/python3 "$HOME/.agentic-sniper/run_paper.py" --min-wallet 2 --pages 5 --max-trades 10 &gt;&gt; "$HOME/.agentic-sniper/paper.log" 2&gt;&amp;1</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- v0.3.1: launchd does not load .zshrc; explicitly forward proxy + auth env -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HTTP_PROXY</key>
        <string>http://127.0.0.1:1087</string>
        <key>HTTPS_PROXY</key>
        <string>http://127.0.0.1:1087</string>
        <key>http_proxy</key>
        <string>http://127.0.0.1:1087</string>
        <key>https_proxy</key>
        <string>http://127.0.0.1:1087</string>
        <key>NO_PROXY</key>
        <string>localhost,127.0.0.1,::1,.local</string>
        <key>HOME</key>
        <string>/Users/l</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/agentic-sniper-paper.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/agentic-sniper-paper.err</string>
</dict>
</plist>
